# base/src/is_web/themes/silicom/root/etc/uci-defaults/luci-theme-silicom
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Bootstrap=/luci-static/silicom
	commit luci
EOF
exit 0
