# product/ubmc/src/ubmc_web/themes/att/root/etc/uci-defaults/luci-theme-att
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Bootstrap=/luci-static/att
	commit luci
EOF
exit 0
